# crontab
0  9 * * * cd ~/.dotfiles && /opt/homebrew/bin/brew bundle dump --force && make format
0  9 * * 2 cd ~/.config/nvim && git fetch && git rebase origin/master && nvim --headless "+Lazy! restore" +qa
0 10 * * * /opt/homebrew/bin/brew update && /opt/homebrew/bin/brew upgrade
